# scripts/fizzbuzz_server.service
[Unit]
Description=Fizz Buzz Server
[Service]
Type=simple
Restart=always
RestartSec=1
ExecStart=go 

[Install]
WantedBy=multi-user.target